# etc/profile.d/inputrc.sh
if [ ! -f ~/.inputrc ]; then
	export INPUTRC="/etc/inputrc"
fi

