# etc/profile.d/ls.sh
eval $(dircolors)

alias ls='ls --color=auto'
